# sys-process/scron/files/init-9999
#!/sbin/runscript

command=crond
name="cron daemon"
